# com.deepseek.balance.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.deepseek.balance</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/swift</string>
        <string>/Users/nick/Documents/DeepSeek 鵝숅쥫/DeepSeekBalance.swift</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/nick/Documents/DeepSeek 鵝숅쥫</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/nick/Documents/DeepSeek 鵝숅쥫/app.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/nick/Documents/DeepSeek 鵝숅쥫/app.err.log</string>
</dict>
</plist>
